# conf.d/webkvm.conf
[program:webkvm]
command=bash /app/webkvm/startup.sh
autorestart=false
